# testbed/isp-lxc/kea-rebind.service
[Unit]
Description=Re-establish eth0-bound sim services (kea, pd-route) when eth0 is recreated
After=network-online.target kea-dhcp6-server.service pd-route.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/kea-rebind
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
