# testbed/isp-lxc/pd-route.service.tmpl
[Unit]
Description=Route delegated /60 prefix back to VM 950
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/ip -6 route replace __PD_PREFIX__ via __VM950_LL__ dev eth0
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
