# testbed/isp-lxc/pd-route.service.j2
[Unit]
Description=Route delegated prefix back to VM 950
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/ip -6 route replace {{ isp.pd_prefix }} via {{ isp.vm950_ll }} dev eth0
{% if isp.static_v4_block %}
# Routed static /29: deliver the block to the WAN VM's stable link address so
# inbound 1:1-NAT'd traffic reaches it (models AT&T routing a static block).
ExecStart=/sbin/ip route replace {{ isp.static_v4_block }} via {{ isp.static_v4_route_to }} dev eth0
{% endif %}
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
