# testbed/isp-lxc/pd-route.service.j2
[Unit]
Description=Route delegated prefix back to VM 950
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/ip -6 route replace {{ isp.pd_prefix }} via {{ isp.vm950_ll }} dev eth0
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
