# huly/timers/huly-update-checker.timer
[Unit]
Description=Run huly-update-checker weekly

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
