# testbed/isp-lxc/200/pd-route.service
[Unit]
Description=Route delegated prefix back to VM 950
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/ip -6 route replace 3d06:bad:b01:2200::/56 via fe80::be24:11ff:febe:8eb4 dev eth0
# Routed static /29: deliver the block to the WAN VM's stable link address so
# inbound 1:1-NAT'd traffic reaches it (models AT&T routing a static block).
ExecStart=/sbin/ip route replace 10.241.204.0/29 via 10.240.204.2 dev eth0
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
