# testbed/isp-lxc/202/pd-route.service
[Unit]
Description=Route delegated prefix back to VM 950
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/ip -6 route replace 3d06:bad:b01:2400::/56 via fe80::be24:11ff:fe3d:cecc dev eth0
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
