# mwan/go/cmd/mwan/opnsense-src/etc/rc.d/mwan_opnsense
#!/bin/sh
# shellcheck shell=sh disable=SC2034
# PROVIDE: mwan_opnsense
# REQUIRE: NETWORKING FILESYSTEMS
# KEYWORD: shutdown
#
# The mwan opnsense daemon takes no flags; its config lives in
# /var/lib/mwan/daemon.toml (MWAN-193) or compiled defaults today.
# Enable in /etc/rc.conf with:
#   mwan_opnsense_enable="YES"

. /etc/rc.subr

name="mwan_opnsense"
rcvar="${name}_enable"

load_rc_config "$name"

: ${mwan_opnsense_enable:="NO"}
: ${mwan_opnsense_logfile:="/var/log/mwan-opnsense.log"}
: ${mwan_opnsense_stop_timeout:="10"}

pidfile="/var/run/${name}.pid"
command="/usr/local/sbin/mwan-opnsense"
required_files="${command}"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

mwan_opnsense_status() {
	if [ -f "${pidfile}" ]; then
		pid=$(cat "${pidfile}" 2>/dev/null)
		if [ -n "${pid}" ] && kill -0 "${pid}" 2>/dev/null; then
			echo "${name} is running as pid ${pid}"
			return 0
		fi
	fi
	echo "${name} is not running"
	return 1
}

mwan_opnsense_preflight() {
	# Self-deploy revert path. If the previous start was a fresh deploy
	# (pending-verify marker present) and the daemon never reported
	# health=ok, swap .current with .previous before this start.
	pending="/var/run/mwan-opnsense.pending-verify"
	state="/var/run/mwan-opnsense.deploy.state"
	sbin_dir="/usr/local/sbin"

	[ -f "${pending}" ] || return 0

	if awk -F: '$1=="health" && $2=="ok"' "${state}" 2>/dev/null | grep -q .; then
		rm -f "${pending}"
		return 0
	fi

	if [ ! -f "${sbin_dir}/mwan-opnsense.previous" ]; then
		echo "${name}: pending-verify present but no .previous to revert to" >&2
		rm -f "${pending}"
		return 0
	fi

	logger -t mwan-opnsense-deploy "self-check failed; reverting to .previous"
	cp -f "${sbin_dir}/mwan-opnsense.previous" "${sbin_dir}/mwan-opnsense.current" \
		|| { echo "${name}: revert copy failed" >&2; return 1; }
	rm -f "${pending}"
}

mwan_opnsense_start() {
	if mwan_opnsense_status >/dev/null 2>&1; then
		echo "${name} already running"
		return 0
	fi

	mwan_opnsense_preflight || return 1

	rm -f "${pidfile}"
	# daemon(8) backgrounds the foreground binary, owns the pidfile,
	# and redirects stdout/stderr to ${mwan_opnsense_logfile}.
	/usr/sbin/daemon -P "${pidfile}" -o "${mwan_opnsense_logfile}" "${command}"

	count=0
	while [ "${count}" -lt 10 ]; do
		if mwan_opnsense_status >/dev/null 2>&1; then
			mwan_opnsense_status
			return 0
		fi
		count=$((count + 1))
		sleep 1
	done

	echo "${name} failed to start" >&2
	return 1
}

mwan_opnsense_stop() {
	pid=$(cat "${pidfile}" 2>/dev/null) || {
		echo "${name} not running"
		return 0
	}

	case "${pid}" in
		*[!0-9]* | "")
			echo "${name} has invalid pidfile ${pidfile}" >&2
			return 1
			;;
	esac

	if ! kill -0 "${pid}" 2>/dev/null; then
		echo "${name} not running; removing stale pidfile"
		rm -f "${pidfile}"
		return 0
	fi

	kill -TERM "${pid}" 2>/dev/null || true
	count=0
	while [ "${count}" -lt "${mwan_opnsense_stop_timeout}" ]; do
		if ! kill -0 "${pid}" 2>/dev/null; then
			rm -f "${pidfile}"
			echo "${name} stopped"
			return 0
		fi
		count=$((count + 1))
		sleep 1
	done

	echo "${name} did not stop after TERM; sending KILL" >&2
	kill -KILL "${pid}" 2>/dev/null || true
	rm -f "${pidfile}"
}

run_rc_command "$1"
