# mwan/go/cmd/mwan/opnsense-src/etc/rc.d/mwan_opnsense
#!/bin/sh
# shellcheck shell=sh disable=SC2034,SC2223,SC1091
# SC2034: rcvar/cmd vars are consumed by rc.subr, not this file.
# SC2223: the `: ${var:=default}` rc.conf idiom assigns static literals.
# SC1091: /etc/rc.subr is a FreeBSD system file, not resolvable here.
# PROVIDE: mwan_opnsense
# REQUIRE: NETWORKING FILESYSTEMS
# KEYWORD: shutdown
#
# The mwan opnsense daemon takes no flags; its runtime config lives in
# /var/lib/mwan/daemon.toml. This rc.d script templates that TOML from
# rc.conf.d-overridable variables before starting the daemon, so the
# rc.conf knobs below are the operator's surface and the TOML is the
# runtime contract the daemon reads.
#
# Enable in /etc/rc.conf with:
#   mwan_opnsense_enable="YES"
#
# Operator knobs (defaults shown):
#   mwan_opnsense_listen_serial="/dev/ttyV0.1"
#   mwan_opnsense_baud="115200"
#   mwan_opnsense_config_xml="/conf/config.xml"
#   mwan_opnsense_backup_dir="/conf/backup"
#   mwan_opnsense_logfile="/var/log/mwan-opnsense.log"
#   mwan_opnsense_state_dir="/var/lib/mwan/transfers"

. /etc/rc.subr

name="mwan_opnsense"
rcvar="${name}_enable"

load_rc_config "$name"

: ${mwan_opnsense_enable:="NO"}
: ${mwan_opnsense_listen_serial:="/dev/ttyV0.1"}
: ${mwan_opnsense_baud:="115200"}
: ${mwan_opnsense_config_xml:="/conf/config.xml"}
: ${mwan_opnsense_backup_dir:="/conf/backup"}
: ${mwan_opnsense_logfile:="/var/log/mwan-opnsense.log"}
: ${mwan_opnsense_state_dir:="/var/lib/mwan/transfers"}
: ${mwan_opnsense_stop_timeout:="10"}

daemon_toml_dir="/var/lib/mwan"
daemon_toml_path="${daemon_toml_dir}/daemon.toml"

pidfile="/var/run/${name}.pid"
child_pidfile="/var/run/${name}.child.pid"
command="/usr/local/sbin/mwan-opnsense"
run_shim="/usr/local/libexec/mwan-opnsense-run"
required_files="${command} ${run_shim}"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

mwan_opnsense_status() {
	if [ -f "${pidfile}" ]; then
		pid=$(cat "${pidfile}" 2>/dev/null)
		if [ -n "${pid}" ] && kill -0 "${pid}" 2>/dev/null; then
			echo "${name} is running as pid ${pid}"
			return 0
		fi
	fi
	echo "${name} is not running"
	return 1
}

mwan_opnsense_write_daemon_toml() {
	# Template /var/lib/mwan/daemon.toml from the rc.conf.d-overridable
	# variables above. The daemon reads this file at startup and refuses
	# to fall back to compiled defaults, so the file must exist with
	# every required field before exec'ing the daemon.
	if [ ! -d "${daemon_toml_dir}" ]; then
		mkdir -p "${daemon_toml_dir}" || {
			echo "${name}: mkdir ${daemon_toml_dir} failed" >&2
			return 1
		}
		chmod 700 "${daemon_toml_dir}" || {
			echo "${name}: chmod 700 ${daemon_toml_dir} failed" >&2
			return 1
		}
	fi

	tmp="${daemon_toml_path}.tmp.$$"
	umask 077
	cat > "${tmp}" <<EOF
# Auto-generated by /usr/local/etc/rc.d/${name}. Do not hand-edit; the
# operator's knobs live in /etc/rc.conf.d/${name} (or /etc/rc.conf).
[daemon]
serial_path = "${mwan_opnsense_listen_serial}"
baud = ${mwan_opnsense_baud}
config_xml_path = "${mwan_opnsense_config_xml}"
backup_dir = "${mwan_opnsense_backup_dir}"
logfile = "${mwan_opnsense_logfile}"
state_dir = "${mwan_opnsense_state_dir}"
EOF
	chmod 600 "${tmp}" || {
		echo "${name}: chmod 600 ${tmp} failed" >&2
		rm -f "${tmp}"
		return 1
	}
	mv -f "${tmp}" "${daemon_toml_path}" || {
		echo "${name}: install ${daemon_toml_path} failed" >&2
		rm -f "${tmp}"
		return 1
	}
}

mwan_opnsense_start() {
	if mwan_opnsense_status >/dev/null 2>&1; then
		echo "${name} already running"
		return 0
	fi

	mwan_opnsense_write_daemon_toml || return 1

	rm -f "${pidfile}" "${child_pidfile}"
	# daemon(8) supervises and restarts the child on an unsolicited exit
	# (-r), records the supervisor pid (-P) and the child pid (-p) so a
	# forced stop can signal the actual daemon, and redirects stdout/stderr
	# to ${mwan_opnsense_logfile}. The supervised target is the preflight
	# shim, not the daemon directly: daemon(8) re-execs the shim on every
	# (re)spawn, so the deploy preflight (revert a bad self-deploy to
	# .previous) runs before the first start and before every respawn. A
	# SIGTERM to the supervisor (service stop) forwards to the child and
	# exits without respawning, so stop still means down.
	/usr/sbin/daemon -r -P "${pidfile}" -p "${child_pidfile}" -o "${mwan_opnsense_logfile}" "${run_shim}"

	count=0
	while [ "${count}" -lt 10 ]; do
		if mwan_opnsense_status >/dev/null 2>&1; then
			mwan_opnsense_status
			return 0
		fi
		count=$((count + 1))
		sleep 1
	done

	echo "${name} failed to start" >&2
	return 1
}

mwan_opnsense_stop() {
	pid=$(cat "${pidfile}" 2>/dev/null) || {
		echo "${name} not running"
		return 0
	}

	case "${pid}" in
		*[!0-9]* | "")
			echo "${name} has invalid pidfile ${pidfile}" >&2
			return 1
			;;
	esac

	if ! kill -0 "${pid}" 2>/dev/null; then
		echo "${name} not running; removing stale pidfile"
		rm -f "${pidfile}" "${child_pidfile}"
		return 0
	fi

	# TERM the daemon(8) supervisor; it forwards SIGTERM to the child and
	# waits for the child to exit.
	kill -TERM "${pid}" 2>/dev/null || true
	count=0
	while [ "${count}" -lt "${mwan_opnsense_stop_timeout}" ]; do
		if ! kill -0 "${pid}" 2>/dev/null; then
			rm -f "${pidfile}" "${child_pidfile}"
			echo "${name} stopped"
			return 0
		fi
		count=$((count + 1))
		sleep 1
	done

	# Forced path. A child wedged in a blocking serial read does not exit
	# on SIGTERM, and a SIGKILL aimed only at the supervisor cannot be
	# forwarded, so the wedged child would orphan holding the serial
	# device. The order is load-bearing: kill the supervisor FIRST, so that
	# under daemon(8) -r (added later) it stops supervising and cannot
	# respawn the child we are about to kill. Then signal the supervisor's
	# process group: daemon(8) calls setsid, so the supervisor leads the
	# group the child shares, and the group id outlives the dead leader, so
	# a negative pid still reaches the wedged child. Finally signal the
	# tracked child directly as a fallback in case it sits in its own group.
	echo "${name} did not stop after TERM; sending KILL" >&2
	kill -KILL "${pid}" 2>/dev/null || true
	kill -KILL -- "-${pid}" 2>/dev/null || true
	child_pid=$(cat "${child_pidfile}" 2>/dev/null)
	case "${child_pid}" in
		*[!0-9]* | "") child_pid="" ;;
	esac
	if [ -n "${child_pid}" ]; then
		kill -KILL "${child_pid}" 2>/dev/null || true
	fi
	rm -f "${pidfile}" "${child_pidfile}"
}

run_rc_command "$1"
