# testbed/isp-lxc/201/pd-route.service
[Unit]
Description=Route delegated prefix back to VM 950
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/ip -6 route replace 3d06:bad:b01:2300::/60 via fe80::be24:11ff:fec0:d760 dev eth0
# Routed static /29: deliver the block to the WAN VM's stable link address so
# inbound 1:1-NAT'd traffic reaches it (models AT&T routing a static block).
ExecStart=/sbin/ip route replace 10.241.205.0/29 via 10.240.205.2 dev eth0
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
