# mwan/go/cmd/mwan/mwan-oob.service
[Unit]
Description=MWAN OOB management daemon (mbrains addrs/rules/routes, DHCP)
Documentation=file:///root/configs/mwan/CUTOVER2-RUNBOOK.md
After=networking.service
Wants=networking.service

[Service]
Type=simple
ExecStart=/usr/local/bin/mwan oob
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mwan-oob

# The daemon needs:
#   CAP_NET_ADMIN  - ip addr/rule/route
#   CAP_NET_RAW    - DHCPv4 raw sockets (insomniacslk/dhcp/nclient4)
#   CAP_NET_BIND_SERVICE - allows bind to low ports if ever needed
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

# Lock down the rest of the filesystem.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
PrivateTmp=true
ReadWritePaths=/var/log /var/run /run

[Install]
WantedBy=multi-user.target
