# proxmox/systemd/wan-monitor.service
[Unit]
Description=WAN Failover Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/scripts/wan-monitor
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

