# mwan/services/mwan-change-detect.service
[Unit]
Description=MWAN config change detector

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mwan-change-detect
