# mwan/services/mwan-update-npt.service
[Unit]
Description=Reapply MWAN IPv6 NPT rules at boot
After=network-online.target nftables.service
Wants=network-online.target nftables.service

[Service]
Type=oneshot
# If we start too early (interfaces not present yet, nftables not ready), retry until it works.
Restart=on-failure
RestartSec=5

EnvironmentFile=/etc/mwan/mwan.env

ExecStartPre=/usr/local/bin/mwan-wait-npt-prereqs.sh

ExecStart=/usr/local/bin/mwan-update-npt-all.sh

[Install]
WantedBy=multi-user.target

