# powerdns/powerdns-admin.service.j2
[Unit]
Description=PowerDNS-Admin Web Interface
After=network.target postgresql.service pdns.service
Requires=postgresql.service

[Service]
Type=simple
User={{ powerdns_admin_user }}
Group={{ powerdns_admin_group }}
WorkingDirectory={{ powerdns_admin_home }}
Environment="PATH={{ powerdns_admin_home }}/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_APP={{ powerdns_admin_home }}/powerdnsadmin/__init__.py"
ExecStart={{ powerdns_admin_home }}/venv/bin/gunicorn --bind 127.0.0.1:{{ powerdns_admin_port }} --workers 2 --threads 2 --timeout 120 --access-logfile {{ powerdns_admin_home }}/logs/access.log --error-logfile {{ powerdns_admin_home }}/logs/error.log powerdnsadmin:app
Restart=on-failure
RestartSec=5s

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ powerdns_admin_home }}/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=powerdns-admin

[Install]
WantedBy=multi-user.target

