# proxmox/services/mwan-watchdog.service.j2
[Unit]
Description=MWAN Auto-Rollback Watchdog (Go)
After=network.target
Wants=network.target

[Service]
Type=simple
EnvironmentFile=/etc/mwan-watchdog/watchdog.env
ExecStart=/usr/local/bin/mwan-watchdog
Restart=always
RestartSec=10
User=root
Group=root
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
